# ios-pcie/layout/Library/LaunchDaemons/com.ttdennis.internalblued.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ttdennis.internalblued</string>
    <key>MachServices</key>
    <dict>
        <key>com.ttdennis.internalblued</key>
        <true/>
    </dict>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/internalblued</string>
    </array>
    <key>UserName</key>
    <string>mobile</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/var/mobile/internalblued.log</string>
    <key>StandardOutPath</key>
    <string>/var/mobile/internalblued.log</string>
</dict>
</plist>
